# paquimetro/crontab/notification_receipt.sh
#!/bin/sh

curl -X 'POST' \
  'http://localhost:8081/notification/receipt' \
  -H 'accept: */*' \
  -d ''